# data/supervisord.conf
[supervisord]
logfile=/var/log/owncloud/supervisord.log ; (main log file;default $CWD/supervisord.log)

[program:cron]
command = /usr/sbin/cron -f -L 15
autostart=true
autorestart=true
stdout_logfile=/var/log/owncloud/supervisord-%(program_name)s.log
stderr_logfile=/var/log/owncloud/supervisord-%(program_name)s.log

[program:phpfpm]
command = /usr/sbin/php5-fpm -F -c /etc/php5/fpm
stdout_logfile=/var/log/owncloud/supervisord-%(program_name)s.log
stderr_logfile=/var/log/owncloud/supervisord-%(program_name)s.log

[program:nginx]
command = /usr/sbin/nginx
stdout_logfile=/var/log/owncloud/supervisord-%(program_name)s.log
stderr_logfile=/var/log/owncloud/supervisord-%(program_name)s.log
